# trading-bot.service
[Unit]
Description=Lotus Koril Trading Bot Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/lotus-koril
EnvironmentFile=/opt/lotus-koril/.env
ExecStart=/usr/bin/bundle exec ruby /opt/lotus-koril/ea.rb
Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lotus-koril

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/lotus-koril/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
